# srcs/mysql/supervisord.conf
[supervisord]
user=root
nodaemon = true

[program: mysql]
command=/usr/sbin/mysqld
autorestart=true